# deploy/macos/com.hr121.docker.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hr121.docker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/ABSOLUTE/PATH/TO/hr121-web/scripts/docker-start-on-login.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/hr121-web</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/hr121-web/logs/docker-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/hr121-web/logs/docker-launchd.err.log</string>
</dict>
</plist>
